# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/echo-project/supervisord.log
logfile_maxbytes=20MB
logfile_backups=5
loglevel=info
pidfile=/run/supervisord.pid

; ============================================================
; 1. Broker MQTT — démarre en premier (requis par tout le reste)
; ============================================================
[program:mosquitto]
command=/usr/sbin/mosquitto -c /etc/mosquitto/mosquitto.conf
priority=10
autostart=true
autorestart=true
startsecs=2
startretries=5
stdout_logfile=/var/log/echo-project/mosquitto.log
stdout_logfile_maxbytes=5MB
redirect_stderr=true

; ============================================================
; 2. GStreamer Daemon — démarre après mosquitto
; ============================================================
[program:gstd]
command=/usr/local/bin/gstd
priority=20
autostart=true
autorestart=true
startsecs=5
startretries=5
stdout_logfile=/var/log/echo-project/gstd.log
stdout_logfile_maxbytes=5MB
redirect_stderr=true

; ============================================================
; 3. Interface Node-RED
; Lancé depuis /usr/src/node-red pour que les plugins
; installés dans ce node_modules soient automatiquement trouvés.
; /data est monté en volume → flows.json et settings.js du dépôt
; ============================================================
[program:node-red]
command=node /usr/src/node-red/node_modules/.bin/node-red --userDir /data
directory=/usr/src/node-red
priority=30
autostart=true
autorestart=true
startsecs=10
startretries=5
environment=TZ="%(ENV_TZ)s"
stdout_logfile=/var/log/echo-project/node-red.log
stdout_logfile_maxbytes=10MB
redirect_stderr=true

; ============================================================
; 4. Serveur Python ECHO — attend 20s que gstd soit prêt
; ============================================================
[program:echo-server]
command=bash -c "sleep 20 && python3 /app/server.py"
priority=40
autostart=true
autorestart=true
startsecs=25
startretries=5
stdout_logfile=/var/log/echo-project/echo-server.log
stdout_logfile_maxbytes=10MB
redirect_stderr=true
